# deploy/com.atlas.assistant.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- ATLAS LaunchAgent template. Paths are filled in by deploy/install-service.sh.
     KeepAlive keeps ATLAS running (restarts on crash) and lets the self-updater
     restart cleanly: after an update, the process exits and launchd respawns it
     with the new code. ATLAS_SUPERVISED tells the app a supervisor is present. -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.atlas.assistant</string>
  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>-m</string>
    <string>atlas.server</string>
  </array>
  <key>WorkingDirectory</key><string>__REPO__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>ATLAS_SUPERVISED</key><string>1</string>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>5</integer>
  <key>ProcessType</key><string>Interactive</string>
  <key>StandardOutPath</key><string>__REPO__/atlas/logs/launchd.out.log</string>
  <key>StandardErrorPath</key><string>__REPO__/atlas/logs/launchd.err.log</string>
</dict>
</plist>
